# ops/systemd/missioncontrol.service.example
[Unit]
Description=MissionControl Next.js App
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=missioncontrol
Group=missioncontrol
WorkingDirectory=/opt/missioncontrol
Environment=NODE_ENV=production
EnvironmentFile=/opt/missioncontrol/.env
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
